# deploy/supervisord.conf
; Unified process manager for the all-in-one image.
; One container runs: chatgpt2api (public :80) + vn-mcp-hub (internal :8005)
; + captcha-solver API (internal :8010)
; + the headful browser stack (Xvfb / Fluxbox / x11vnc / noVNC :6080)
; used for manual web logins.
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── Headful browser stack (for captcha-solver manual logins via noVNC) ──
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1366x768x24 -ac +extension RANDR
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":99"
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Đặt env VNC_PASSWORD trên container để khóa noVNC/x11vnc bằng mật khẩu;
; không đặt = -nopw như trước (chỉ nên dùng trong LAN, đừng expose 6080 public).
[program:x11vnc]
command=/bin/sh -c 'if [ -n "$VNC_PASSWORD" ]; then exec /usr/bin/x11vnc -display :99 -passwd "$VNC_PASSWORD" -forever -shared -rfbport 5900 -quiet; else exec /usr/bin/x11vnc -display :99 -nopw -forever -shared -rfbport 5900 -quiet; fi'
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ 6080 localhost:5900
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── vn-mcp-hub (internal only, reached via 127.0.0.1:8005) ──
[program:vn-mcp-hub]
command=/app/.venv/bin/uvicorn src.main:app --host 127.0.0.1 --port 8005 --workers 1
directory=/app/mcp_hub
environment=PYTHONPATH="/app/mcp_hub"
autorestart=true
priority=45
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── captcha-solver API (internal only — the web UI reaches it via the
;    same-origin /api/captcha proxy; server-side callers use 127.0.0.1:8010) ──
[program:captcha-api]
command=/app/.venv/bin/uvicorn src.main:app --host 127.0.0.1 --port 8010 --workers 1
directory=/app/captcha
environment=DISPLAY=":99",PYTHONPATH="/app/captcha"
autorestart=true
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── zalo-server (Node zca-js, Zalo cá nhân) ──
;    Nội bộ 127.0.0.1:3001 cho chatgpt2api (services/zalo_personal.py) gọi;
;    cũng bind 0.0.0.0:3001 để Home Assistant integration trỏ thẳng vào.
;    Data (cookies/sessions/messages/users) persist dưới /app/data/zalo_bot.
[program:zalo-server]
command=/usr/local/bin/node server.js
directory=/app/zalo-server
environment=PORT="3001",DATA_DIRECTORY="/app/data/zalo_bot",NODE_ENV="production",TZ="Asia/Ho_Chi_Minh"
autorestart=true
priority=55
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── chatgpt2api (public app + web UI + MCP admin proxy) ──
[program:chatgpt2api]
command=/app/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 80 --access-log
directory=/app
environment=DISPLAY=":99",PYTHONPATH="/app"
autorestart=true
priority=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
